# deploy/supervisor/conf.d/redis.conf
[program:redis]
directory=/home/task2/mysites/
command=/usr/bin/redis-server 
autostart=true
autorestart=true
redirect_stderr=true
