# deploy/supervisor/conf.d/sentry.conf
[program:sentry]
directory=/home/hadoop/.sentry/
command=/usr/local/bin/sentry start http
autostart=true
autorestart=true
redirect_stderr=true
